# caddy2/rootfs/etc/services.d/caddy2/run
#!/usr/bin/with-contenv sh

if [ "$docker-proxy" == "true" ];then
  exec /bin/caddy docker-proxy
else
  exec /bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
fi