# caddy2/rootfs/etc/services.d/caddy2/run
#!/usr/bin/with-contenv sh

exec /bin/caddy docker-proxy