# caddy2/rootfs/etc/services.d/caddy2/run
#!/usr/bin/with-contenv sh

exec caddy docker-proxy