# elecV2P/rootfs/etc/services.d/node/run
#!/usr/bin/with-contenv sh

exec /bin/sh -c 'cd /usr/local/app && pm2-runtime index.js'