# git/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

#启动cron
exec crond -f -L /dev/null